# mpd/root/etc/s6-overlay/s6-rc.d/mpd/run
#!/usr/bin/with-contenv sh

set -e

mountpoint $CACHE_MOUNT_PATH

exec mpd \
  --no-daemon \
  --stdout \
  $MPD_CONF_PATH