# sunshine-desktop/root/etc/s6-overlay/s6-rc.d/svc-udev/run
#!/usr/bin/with-contenv bash

exec /lib/systemd/systemd-udevd