# sunshine-desktop/root/etc/s6-overlay/s6-rc.d/svc-de/run
#!/usr/bin/with-contenv bash

set -xe

cd $HOME
exec s6-setuidgid $USER \
  xfce4-session > /dev/null 2>&1