# sunshine-desktop/root/etc/s6-overlay/s6-rc.d/svc-pulseaudio/run
#!/usr/bin/with-contenv bash

set -xe

exec s6-setuidgid $USER \
  pulseaudio \
    --log-level=0 \
    --log-target=stderr \
    --exit-idle-time=-1 > /dev/null 2>&1