# code-server/root/etc/s6-overlay/s6-rc.d/svc-codeserver/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid $USER code-server \
  --auth=none \
  --bind-addr=0.0.0.0:$CODE_PORT