# clickhouse/root/etc/s6-overlay/s6-rc.d/svc-jfs/run
#!/usr/bin/with-contenv sh

mkdir -p $CLICKHOUSE_DATA_PATH

exec juicefs mount \
  --storage minio \
  --bucket $JFS_MINIO_BUCKET \
  sqlite3://$JFS_DB_PATH \
  $CLICKHOUSE_DATA_PATH