# clickhouse/root/etc/s6-overlay/s6-rc.d/svc-clickhouse/run
#!/usr/bin/with-contenv sh

set -e

mountpoint $CLICKHOUSE_DATA_PATH

exec clickhouse-server \
  -C $CLICKHOUSE_CONF_PATH