# code-server/root/etc/s6-overlay/s6-rc.d/svc-jfs/run
#!/usr/bin/with-contenv sh

set -e

mkdir -p $(dirname $HOME)

exec juicefs mount \
  --storage minio \
  --bucket $JFS_MINIO_BUCKET \
  sqlite3://$JFS_DB_PATH \
  $(dirname $HOME)