# sunshine-desktop/root/etc/s6-overlay/s6-rc.d/svc-seatd/run
#!/usr/bin/with-contenv bash

exec /usr/bin/seatd -u abc -g abc