# transmission/root/etc/s6-overlay/s6-rc.d/svc-transmission/run
#!/usr/bin/with-contenv sh

set -e

mountpoint $HOME

cp $TRANSMISSION_CONF_PATH $HOME/settings.json

exec transmission-daemon \
  --foreground \
  --config-dir $HOME