# clickhouse/root/etc/s6-overlay/s6-rc.d/init-jfs/run
#!/usr/bin/with-contenv sh

juicefs format \
  --storage minio \
  --bucket $JFS_MINIO_BUCKET \
  --access-key $JFS_MINIO_ACCESS_KEY_ID \
  --secret-key $JFS_MINIO_SECRET_ACCESS_KEY \
  sqlite3://$JFS_DB_PATH \
  $JFS_RESOURCE_NAME