# kasm-desktop/root/etc/s6-overlay/s6-rc.d/init-user/run
#!/usr/bin/with-contenv bash

set -xe

mkdir -p $HOME
chown $UID:$UID $HOME
useradd $USER -d $HOME -m -u $UID
usermod -G wheel,video $USER

exec s6-setuidgid $USER \
  cp -r /etc/skel/. $HOME