# sunshine-desktop/root/etc/s6-overlay/s6-rc.d/init-de/run
#!/usr/bin/with-contenv bash

set -xe

rm -f \
  /etc/xdg/autostart/xfce-polkit.desktop \
  /etc/xdg/autostart/at-spi-dbus-bus.desktop